# dummy.service
[Unit]
Description=Dummy Background Service demonstrating systemd management
# Start after the network is up, adjust if needed
After=network.target

[Service]
# Type=simple means the ExecStart command is the main process
Type=simple

# Specify the user and group to run as (optional but good practice)
# Create the user first: sudo useradd -r -s /bin/false dummyuser
# User=dummyuser
# Group=dummyuser
# If no User/Group specified, it defaults to root

# Full path to the executable script
ExecStart=/usr/local/bin/dummy.sh

# Restart the service alwasyif it exits for an reason (crash, munual stop followed by unexpected exit)
Restart=failure
# Wait 5 seconds before attempting restart
RestartSec=5s

# Redirect standard output and standad error to the systemd journal
StandardOutput=journal
StandardError=journal

# Optional: Set a working directory if the script needs it
# WorkingDirectory=/path/to/working/dir

[Install]
# Enable the service for the default multi-user target (typical runlevel)
WantedBy=multi-user.target

